# build/container/etc/s6-overlay/s6-rc.d/pbs-plus-agent-log-prepare/up
if { mkdir -p /var/log/pbsplus }
if { chown nobody:nogroup /var/log/pbsplus }
chmod 02755 /var/log/pbsplus
