# build/container/etc/s6-overlay/s6-rc.d/pbs-plus-agent/run
#!/command/with-contenv sh
set -eu

exec 2>&1

# Environment
HOME="/var/lib/pbs-plus-agent"; export HOME
USER="pbsplus"; export USER

# Create required dirs with ownership
mkdir -p /run/pbs-plus-agent \
         /var/lib/pbs-plus-agent \
         /var/log/pbs-plus-agent \
         /etc/pbs-plus-agent

chown pbsplus:pbsplus /run/pbs-plus-agent \
                      /var/lib/pbs-plus-agent \
                      /var/log/pbs-plus-agent \
                      /etc/pbs-plus-agent

exec /usr/bin/pbs-plus-agent
