# build/package/linux_agent/openrc/pbs-plus-agent
#!/sbin/openrc-run

name="PBS Plus Agent"
description="PBS Plus Agent for backup operations"

command="/usr/bin/pbs-plus-agent"
command_user="pbsplus:pbsplus"
pidfile="/run/pbs-plus-agent/pbs-plus-agent.pid"
directory="/var/lib/pbs-plus-agent"
agent_env="/etc/pbs-plus-agent/agent.env"

output_log="/var/log/pbs-plus-agent/agent.log"
error_log="/var/log/pbs-plus-agent/agent.error.log"

# Background execution & Resource management
# --nicelevel 19: Lowest CPU priority
# --ionice idle: Only perform I/O when system is idle (prevents disk lag)
command_background="yes"
start_stop_daemon_args="--nicelevel 19 --ionice idle"

depend() {
    need net
    after firewall
}

start_pre() {
    checkpath --directory --owner pbsplus:pbsplus --mode 0755 /run/pbs-plus-agent
    checkpath --directory --owner pbsplus:pbsplus --mode 0750 /var/log/pbs-plus-agent
    checkpath --directory --owner pbsplus:pbsplus --mode 0750 /var/lib/pbs-plus-agent
    
    if [ -f "$agent_env" ]; then
        while IFS= read -r line || [ -n "$line" ]; do
            case "$line" in
                ''|#*) continue ;;
                *=*) 
                    key=$(echo "${line%%=*}" | tr -d ' ')
                    val=$(echo "${line#*=}" | sed "s/^['\"]//;s/['\"]$//")
                    export "$key=$val"
                    ;;
            esac
        done < "$agent_env"
    fi

    # OpenRC doesn't manage caps, so we verify the install script did its job
    if command -v getcap >/dev/null 2>&1; then
        if ! getcap "$command" | grep -q "cap_dac_read_search"; then
            ewarn "Warning: Binary lacks CAP_DAC_READ_SEARCH. Backups may fail to read root files."
        fi
    fi
}

reload() {
    ebegin "Reloading ${name}"
    start-stop-daemon --signal HUP --pidfile "${pidfile}"
    eend $?
}

stop_post() {
    [ -f "${pidfile}" ] && rm -f "${pidfile}"
}
