# build/package/linux_agent/systemd/pbs-plus-agent.service
[Unit]
Description=PBS Plus Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pbsplus
Group=pbsplus
WorkingDirectory=/var/lib/pbs-plus-agent
ExecStart=/usr/bin/pbs-plus-agent
Restart=on-failure
EnvironmentFile=-/etc/pbs-plus-agent/agent.env

# --- Permissions for Backup & Autoupdate ---
# CAP_DAC_READ_SEARCH: Read any file for backups
# CAP_DAC_OVERRIDE: Replace its own binary in /usr/bin during update
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE
AmbientCapabilities=CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE
NoNewPrivileges=yes

# --- Resource Management ---
Nice=19
IOSchedulingClass=idle

# --- Sandboxing & Write Access ---
# ProtectSystem=full makes /usr and /etc read-only by default
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=yes

ReadWritePaths=/var/lib/pbs-plus-agent /var/log/pbs-plus-agent /run/pbs-plus-agent /etc/pbs-plus-agent /tmp /var/tmp /usr/bin /etc/systemd/system

[Install]
WantedBy=multi-user.target
